# recipes-core/sys-setup/sys-setup/scripts/volatiles
#!/sbin/openrc-run

description="Create volatile directories."

killproc() {
    PID=`pidof $1`
    [ -n "$PID" ] && kill $PID
}

depend()
{
    before bootmisc
    need localmount
    provide /tmp /var/log
}

start() {
    if [ "$VOLATILE_TMP_DIR" == "yes" ];
    then
        mkdir -p /var/volatile/tmp
    fi

    if [ "$VOLATILE_LOG_DIR" == "yes" ];
    then
        mkdir -p /var/volatile/log
    fi

    /etc/init.d/populate-volatile.sh
}
