# recipes-core/sys-setup/sys-setup/profile.d/sysenv.sh
uid=`id -u`
mkdir -p /var/run/$uid

export XDG_RUNTIME_DIR=/var/run/$uid
