# recipes-connectivity/bluez5/bluez5/init
#!/sbin/openrc-run

description="bluetooth."

depend()
{
    need dbus
    provide bluetooth bluz5
}

command="$DAEMON"
command_args="$NOPLUGIN_OPTION"
command_background=true
pidfile="${PIDFILE}"
supervisor=supervise-daemon
